# old_projects_watchdog/systemd_services/old_projects_watchdog.timer
[Unit]
Description=Updating the scraped Pypi list.
Requires=old_projects_watchdog.service

[Timer]
OnActiveSec=2m

[Install]
WantedBy=timers.target

